# shinken_thruk_graphite/supervisor/conf.d/apache2.conf
[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
user=root
log_stdout=true
log_stderr=true
logfile=/var/log/supervisor_apache2.log
