# shinken_thruk_graphite/supervisor/conf.d/carbon-cache.conf
[program:carbon-cache]
command = python /opt/graphite/bin/carbon-cache.py --debug start
autorestart = true
log_stdout=false
log_stderr=true
logfile = /var/log/supervisor_%(program_name)s.log
