# supervisor/supervisord.d/initialize_docker_shinken.conf
[program:initialize_docker_shinken]
command=/bin/initialize_docker_shinken.sh
autostart=true
log_stdout=false
log_stderr=true
