# supervisor/supervisord.d/shinken-reactionner.conf
[program:shinken-reactionner]
command=/usr/bin/shinken-reactionner -c /etc/shinken/daemons/reactionnerd.ini
autostart=true
user=shinken
log_stdout=true
log_stderr=true