# shinken_basic/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
autostart=true
user=root
log_stdout=true
log_stderr=true
logfile=/var/log/supervisor_nginx.log
