# shinken_basic/supervisor/conf.d/ntpd.conf
[program:ntpd]
command=/usr/sbin/ntpd -u ntp:ntp -p /var/run/ntpd.pid -g -n
process_name=%(program_name)s
autostart=true
autorestart=true
user=root
directory=/
priority=958
stdout_logfile=/var/log/supervisor_%(program_name)s-stdout.log
stderr_logfile=/var/log/supervisor_%(program_name)s-stderr.log
