# shinken_basic/supervisor/conf.d/watch_shinken_config.conf
[program:watch_shinken_config]
command=/usr/bin/watch_shinken_config.sh
autostart=true
user=root
log_stdout=true
log_stderr=true
logfile=/var/log/watch_shinken_config.log
